# incubator/JailKit/sources/debian/postrm
#!/bin/sh
#DEBHELPER#
set -e
if [ "$1" = "purge" ] ; then
        update-rc.d jailkit remove >/dev/null
fi
